# timer/token.cron
0 5 8,28 * * sh /root/Dev/ALGO/timer/token.sh

